# api/pmerit-api.service
[Unit]
Description=PMERIT AI Multi-Tier Portal API
After=network.target

[Service]
Type=simple
User=gabriel-ai
WorkingDirectory=/home/gabriel-ai/pmerit-ai-platform/api
Environment=PATH=/home/gabriel-ai/.local/bin
ExecStart=/usr/bin/python -m uvicorn main:app --host 0.0.0.0 --port 9001
Restart=always

[Install]
WantedBy=multi-user.target
